# infra/cloth.service
[Unit]
Description=CLOTH API service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/cloth/current
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/npm run start --workspace=api
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
